# scripts/systemd/crm-fields-sync.timer
[Unit]
Description=Cordys CRM Fields Sync Timer
Documentation=https://github.com/1Panel-dev/CordysCRM-skills

[Timer]
OnCalendar=Sun *-*-* 02:00:00
Persistent=true
Unit=crm-fields-sync.service

[Install]
WantedBy=timers.target
